# etc/cron.daily/archive_data
#!/bin/bash
python3 /home/gpstracker/vigor22/archive_data.py
